# rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash
echo "[services.d] started nginx"
# Run nginx processes as picocms:web
nginx
